# redispwned.service
[Unit]
Description=RedisPwned service
ConditionPathExists=/var/app/redispwned
After=network.target redis.service

[Service]
Type=simple
User=redis
Group=redis
LimitNOFILE=1024

Restart=on-failure
RestartSec=10

Environment="GIN_MODE=release"
Environment="PORT=443"
AmbientCapabilities=CAP_NET_BIND_SERVICE

WorkingDirectory=/var/app/redispwned
ExecStart=/usr/local/bin/redispwned

PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/log/redispwned
ExecStartPre=/bin/chown syslog:adm /var/log/redispwned
ExecStartPre=/bin/chmod 755 /var/log/redispwned
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=redispwned

[Install]
WantedBy=multi-user.target
